# roles/ansible-haProxy/templates/haproxy-marathon-bridge.cron.j2
* * * * * root python /usr/local/bin/haproxy-marathon-bridge-dir/haproxy-marathon-bridge.py updateConfig $(cat /etc/haproxy-marathon-bridge/marathons) 2>&1 1>>/tmp/haproxy-marathon-bridge.cron.log